# y.conf
[program:y]
command=python3 /home/curator2/y/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/y/test.err.log
stdout_logfile=/home/curator2/y/test.out.log